# modules/cookery/templates/cookery.conf.erb
description "cookery thin webserver"

start on runlevel [2345]
stop on runlevel [!2345]

exec bundle exec thin start -C <%= thin_config %>/<%= app %>.yml
respawn
